# infra/paganini-daemons.timer
[Unit]
Description=PAGANINI AIOS — Trigger daemon run every 6 hours
Requires=paganini-daemons.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
